# deploy/systemd/olkalou-worker@.service
[Unit]
Description=Ol Kalou live engine worker %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=olkalou
Group=olkalou
WorkingDirectory=/opt/olkalou-live-engine
EnvironmentFile=/etc/olkalou/engine.env
Environment=WORKER_ID=%i
ExecStart=/opt/olkalou-live-engine/.venv/bin/python -m olkalou_engine.cli --root /opt/olkalou-live-engine worker
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/olkalou-live-engine/data

[Install]
WantedBy=multi-user.target
